# deploy/vm01/my-bar-whisky-sync.service
[Unit]
Description=my-bar 주류 시트 양방향 동기화 (fullSync)
After=my-bar-next.service

[Service]
Type=oneshot
ExecStart=/usr/bin/curl -s -m 60 -X POST http://127.0.0.1:3003/my-bar/api/whisky/sync
