# deploy/vm01/my-hobby-whisky-sync.timer
[Unit]
Description=my-hobby 위스키 시트 동기화 주기 실행(10분)

[Timer]
OnBootSec=3min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
